# deploy/systemd/smartheaven-v2.service
[Unit]
Description=Smart Heaven v2 - Home Automation System
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/pedro/smart-heaven-v2
ExecStart=/usr/bin/docker-compose -f docker-compose.prod.yml up -d
ExecStop=/usr/bin/docker-compose -f docker-compose.prod.yml down
User=pedro
Group=pedro

# Restart policy
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
